# Cascade Windsurf/TRADE IT/TradeHub/com.trading.ibkr-service.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.trading.ibkr-service</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/michaelb.williams/Cascade Windsurf/TRADE IT/TradeHub/start_ibkr_service.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/Users/michaelb.williams/Cascade Windsurf/TRADE IT/TradeHub/logs/ibkr_service_launchd_err.log</string>
    <key>StandardOutPath</key>
    <string>/Users/michaelb.williams/Cascade Windsurf/TRADE IT/TradeHub/logs/ibkr_service_launchd_out.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/michaelb.williams/Cascade Windsurf/TRADE IT/TradeHub</string>
</dict>
</plist>
